# kpi/systemd/kpi-heartbeat.service
[Unit]
Description=KPI metrics-vault — hourly heartbeat oneshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/aiagent/assistant/git/history-production/kpi
Environment=KPI_DB=/home/aiagent/assistant/state/kpi.sqlite
Environment=KPI_LOG_DIR=/home/aiagent/assistant/state/kpi-logs
ExecStart=/home/aiagent/assistant/git/history-production/kpi/scripts/heartbeat.sh
StandardOutput=append:/home/aiagent/assistant/state/kpi-logs/heartbeat-systemd.log
StandardError=append:/home/aiagent/assistant/state/kpi-logs/heartbeat-systemd.log

[Install]
WantedBy=default.target
