# systemd/dashboard-kpi-refresh.service
[Unit]
Description=KPI dashboard daily refresh
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/aiagent/assistant/deploys/kpi-dashboard
Environment=PYTHONUNBUFFERED=1
Environment=DASHBOARD_KPI_DB=/home/aiagent/assistant/state/dashboard-kpi.sqlite
ExecStart=/home/aiagent/miniconda3/envs/practicum/bin/python /home/aiagent/assistant/deploys/kpi-dashboard/scripts/run_refresh.py
# StandardOutput/StandardError intentionally NOT set: Python's
# TimedRotatingFileHandler in run_refresh.py owns the log file lifecycle.
# Letting systemd append to the same path would duplicate every line and
# break rotation (Gemini-3.1-pro F-01 R1 HIGH finding).
